# systemd/qsosysd.socket
[Unit]
Description=quicksov privileged platform_profile helper socket

[Socket]
ListenStream=/run/quicksov/qsosysd.sock
SocketUser=root
SocketGroup=quicksov
SocketMode=0660
DirectoryMode=0755
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
