# init/openrc/qsosysd
#!/sbin/openrc-run

description="quicksov privileged platform_profile helper"
command="/usr/bin/qsosysd"
command_background="no"
supervisor="supervise-daemon"
pidfile="/run/${RC_SVCNAME}.pid"
retry="TERM/5/KILL/5"

depend() {
    need localmount
}
